# etc/systemd/system/backup.service
[Unit]
Description="make the backup"
After=network.target

[Service]
ExecStart=/usr/local/bin/backup.sh
